# crontab
@reboot node /home/pi/taskbot/index.js &
